# raspbian-ua-netinst/beamer-clock.service
[Unit]
Description=InfoBeamer Clock

[Service]
Type=oneshot
ExecStart=/srv/info-beamer/clock.py

[Install]
WantedBy=multi-user.target
